# deployment/annotation.service
[Unit]
Description=Gunicorn service for annotation-engine (FastAPI)
After=network.target

[Service]
User=pasta
WorkingDirectory=/home/pasta/annotation/engine

# Pixi-managed gunicorn with uvicorn workers; keep it bound to loopback only.
ExecStart=/home/pasta/annotation/engine/.pixi/envs/default/bin/gunicorn webapp.run:app -w 4 -k uvicorn.workers.UvicornWorker --bind 127.0.0.1:8001

ExecReload=/bin/kill -s HUP $MAINPID
KillMode=mixed
TimeoutStopSec=5
PrivateTmp=true

Restart=always
Type=simple

# Optional: send logs to journald (default behavior). View with:
#   journalctl -u annotation-engine -f

[Install]
WantedBy=multi-user.target